# supervisor.conf
[program:uwsgi-app]
command = uwsgi --ini /home/docker/app/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx -g 'daemon off;'

autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 30
